# supervisor/supervisor.conf
# [program:node]
# command = npm run start:prod:api
[supervisord]

[program:rsyslog]
command = /usr/sbin/rsyslogd -n -iNONE

[program:sshd]
command = /usr/sbin/sshd -D -p 2222
